# snapcast/etc/services.d/avahi/run
#!/usr/bin/with-contenv bashio

#bashio::log.info "Starting Avahi..."
#exec 2>&1
#sleep 5
#exec avahi-daemon --no-drop-root
