# host_iptables/rootfs/etc/s6-overlay/s6-rc.d/init-services/finish
#!/command/with-contenv bashio
# shellcheck shell=bash

STAYUP=$(bashio::config 'Stay_up')
RESTORE=$(bashio::config 'restore')

# Save iptables rules
bashio::log.info "Saving rules"
iptables-save > /data/iptables/save.v4 || true
ip6tables-save > /data/iptables/save.v6 || true
                          
if [ "$STAYUP" == "true" ] && [ "$RESTORE" == "true" ]; then
  for T in filter nat mangle raw security; do
    bashio::log.info "Flushing $T"
    iptables -t $T -F || true
    ip6tables -t $T -X || true
  done
fi
true
