# host_iptables/rootfs/etc/s6-overlay/s6-rc.d/init-services/run
#!/command/with-contenv bashio
# shellcheck shell=bash

bashio::log.info "Setting Custom iptables rules"

STAYUP=$(bashio::config 'Stay_up')
RESTORE=$(bashio::config 'restore')

if [ ! -d /data/iptables/ ]; then
  bashio::log.info "Creating /data/iptables directory"
  mkdir /data/iptables
fi

RULES=/config/rules

IFS=$'\n'
if [ -f $RULES.v4 ]; then
  bashio::log.info "Loading v4 rules"
  cat $RULES.v4 |xargs iptables -v
  # file=$(cat $RULES.v4)
  # for line in $file
  # do
  #   iptables $line || true
  #   bashio::log.info "iptables $line"
  # done
fi

if [ -f $RULES.v6 ]; then
  bashio::log.info "Loading v6 rules"
  cat $RULES.v6 |xargs ip6tables -v
  # file=$(cat $RULES.v6)
  # for line in $file
  # do
  #   ip6tables $line || true
  #   bashio::log.info "ip6tables $line"
  # done
fi
unset IFS

SAVE=/data/iptables/save
if [ "$RESTORE" == "true" ]; then
  if [ -f $SAVE.v4 ]; then
    bashio::log.info "restoring v4"
    iptables-restore $SAVE.v4 || true
  fi
  if [ -f $SAVE.v6 ]; then
    bashio::log.info "restoring v6"
    ip6tables-restore $SAVE.v6 || true
  fi
fi

if [ "$STAYUP" != "true" ]; then
  /run/s6/basedir/bin/halt
else
  true
fi
