# snapcast/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting NGinx..."
exec 2>&1
exec nginx -g 'daemon off;'