# snapcast/etc/services.d/dbus/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting DBUS Daemon..."
exec 2>&1
exec dbus-daemon --nofork --session --print-address